# kiosk/systemd/wrzdj-kiosk.service
# NOTE: This service is installed but NOT enabled. Cage launches from
# /home/kiosk/.bash_profile instead (needs logind seat access for DRM/input).
# This file is kept as a reference for Chromium flags and security hardening.

[Unit]
Description=WrzDJ Kiosk Display
Documentation=https://github.com/thewrz/WrzDJ/tree/main/kiosk
After=wrzdj-wifi-portal.service
Wants=wrzdj-wifi-portal.service

[Service]
Type=simple
User=kiosk
Group=kiosk
EnvironmentFile=/etc/wrzdj-kiosk.conf

# Hide hardware cursor (touchscreen-only)
Environment=WLR_NO_HARDWARE_CURSORS=1
# Software renderer fallback if GPU drivers misbehave
# Uncomment if you see GPU errors in journalctl:
# Environment=WLR_RENDERER=pixman

ExecStart=/usr/bin/cage -- /usr/bin/chromium-browser \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --no-first-run \
    --disable-translate \
    --disable-session-crashed-bubble \
    --disable-component-update \
    --disable-pinch \
    --touch-events=enabled \
    --ozone-platform=wayland \
    --disable-dev-shm-usage \
    --disable-background-networking \
    --disable-sync \
    --metrics-recording-only \
    --disable-default-apps \
    --no-default-browser-check \
    --autoplay-policy=no-user-gesture-required \
    ${EXTRA_CHROMIUM_FLAGS} \
    http://localhost

Restart=on-failure
RestartSec=5

# Prevent runaway memory usage (Chromium can be hungry)
MemoryMax=75%
MemoryHigh=60%

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/home/kiosk

[Install]
WantedBy=multi-user.target
